# templates/clockwork-users.timer.j2
[Unit]
Description=Update users for clockwork.
Requires=clockwork-users.service

[Timer]
Unit=clockwork-users.service
OnCalendar=Sat *-*-* 00:00:00
AccuracySec=1d

[Install]
WantedBy=timers.target
